# deploy/serbiansubtitles.service
[Unit]
Description=Serbian Subtitles web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=serbiansubtitles
Group=serbiansubtitles
WorkingDirectory=/opt/serbiansubtitles
Environment=NODE_ENV=production
EnvironmentFile=/opt/serbiansubtitles/.env
ExecStart=/usr/local/bin/node /opt/serbiansubtitles/server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
UMask=0077
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
